# roles/feederhub/templates/systemd/feederhub-onvif-vips.service.j2
# {{ ansible_managed }}
# Creates one macvlan child of {{ feederhub_onvif_bridge }} per feeder
# ONVIF VIP.  Each child gets its own MAC (set per-feeder in
# defaults/main.yml) so UniFi Protect's adoption flow doesn't dedupe
# both feeders into one entry — Protect reads MACs via ARP at the
# link layer, so synthesising distinct MACs is what makes two feeder
# IPs adopted as two cameras.
#
# Idempotent: addresses + links that already exist are re-added with
# `|| true` so a re-run after a reboot (where systemd-networkd has
# already brought up the parent bridge but our oneshot hasn't yet
# touched the children) doesn't fail.

[Unit]
Description=feederhub ONVIF VIPs (macvlan children of {{ feederhub_onvif_bridge }})
After=network-online.target sys-subsystem-net-devices-{{ feederhub_onvif_bridge }}.device
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Tighten ARP behaviour so each interface only answers ARP requests
# for its own configured IP.  By default the Linux kernel answers
# ARP for any locally-configured IP from any interface, which means
# both macvlan children + the parent bridge would all answer ARP
# for .184 and .185 — collapsing the per-MAC identity we're trying
# to establish.  arp_ignore=1 + arp_announce=2 together give us
# "interface-scoped ARP" semantics for the rest of the lifetime of
# this systemd unit.
ExecStartPre=/usr/sbin/sysctl -w net.ipv4.conf.all.arp_ignore=1
ExecStartPre=/usr/sbin/sysctl -w net.ipv4.conf.all.arp_announce=2
ExecStartPre=/usr/sbin/sysctl -w net.ipv4.conf.default.arp_ignore=1
ExecStartPre=/usr/sbin/sysctl -w net.ipv4.conf.default.arp_announce=2
# Indexed names (onvif0, onvif1, ...) keep us inside Linux's
# IFNAMSIZ limit (15 chars).  The feeder a given index maps to is
# documented by index order in feederhub_onvif_devices.
{% for dev in feederhub_onvif_devices %}
# {{ dev.name }} → onvif{{ loop.index0 }} ({{ dev.mac }}, {{ dev.ip }})
ExecStart=-/usr/sbin/ip link add link {{ feederhub_onvif_bridge }} name onvif{{ loop.index0 }} address {{ dev.mac }} type macvlan mode bridge
ExecStart=-/usr/sbin/ip address add {{ dev.ip }}/24 dev onvif{{ loop.index0 }}
ExecStart=/usr/sbin/ip link set onvif{{ loop.index0 }} up
{% endfor %}
{% for dev in feederhub_onvif_devices %}
ExecStop=-/usr/sbin/ip link delete onvif{{ loop.index0 }}
{% endfor %}

[Install]
WantedBy=multi-user.target
